# docker/kibana/supervisor/kibana.conf
[program:kibana]
command=/usr/share/kibana/bin/kibana --allow-root
autostart=true
