# docker/prod/supervisord.conf
[supervisord]
environment=ALL_JOBS=%(ENV_ALL_JOBS)s,CURRENT_JOB=%(ENV_CURRENT_JOB)s

[program:assistente_ia]
command=/var/www/html/workers-control.sh assistente_ia
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO
startretries=5
startsecs=1

[program:triagem_ia]
command=/var/www/html/workers-control.sh triagem_ia
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO
startretries=5
startsecs=1

[program:indexacao_componente_digital]
command=/var/www/html/workers-control.sh indexacao_componente_digital
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_pessoa]
command=/var/www/html/workers-control.sh indexacao_pessoa
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_modelo]
command=/var/www/html/workers-control.sh indexacao_modelo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:dense_vector]
command=/var/www/html/workers-control.sh dense_vector
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_repositorio]
command=/var/www/html/workers-control.sh indexacao_repositorio
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_processo]
command=/var/www/html/workers-control.sh indexacao_processo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_tese]
command=/var/www/html/workers-control.sh populate_tese
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_tese]
command=/var/www/html/workers-control.sh indexacao_tese
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_componente_digital]
command=/var/www/html/workers-control.sh populate_componente_digital
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_pessoa]
command=/var/www/html/workers-control.sh populate_pessoa
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_modelo]
command=/var/www/html/workers-control.sh populate_modelo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_repositorio]
command=/var/www/html/workers-control.sh populate_repositorio
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_processo]
command=/var/www/html/workers-control.sh populate_processo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:push_count]
command=/var/www/html/workers-control.sh push_count
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:push_mercure]
command=/var/www/html/workers-control.sh push_mercure
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:relatorio_create]
command=/var/www/html/workers-control.sh relatorio_create
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:push_resource]
command=/var/www/html/workers-control.sh push_resource
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:regras_etiqueta]
command=/var/www/html/workers-control.sh regras_etiqueta
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:process_favorito]
command=/var/www/html/workers-control.sh process_favorito
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:download_processo]
command=/var/www/html/workers-control.sh download_processo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:copiar_documentos_processo]
command=/var/www/html/workers-control.sh copiar_documentos_processo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:integracao_gerar_dossie]
command=/var/www/html/workers-control.sh integracao_gerar_dossie
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:integracao_gerar_dossie_tarefa]
command=/var/www/html/workers-control.sh integracao_gerar_dossie_tarefa
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO


[program:envia_ciencia]
command=/var/www/html/workers-control.sh envia_ciencia
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:envia_documento_avulso]
command=/var/www/html/workers-control.sh envia_documento_avulso
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:envia_processo]
command=/var/www/html/workers-control.sh envia_processo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:envia_componentes_digitais]
command=/var/www/html/workers-control.sh envia_componentes_digitais
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:recebe_recibo_de_tramite]
command=/var/www/html/workers-control.sh recebe_recibo_de_tramite
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:recebe_componentes_digitais]
command=/var/www/html/workers-control.sh recebe_componentes_digitais
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:recebe_tramite]
command=/var/www/html/workers-control.sh recebe_tramite
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:sincroniza_barramento]
command=/var/www/html/workers-control.sh sincroniza_barramento
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:tarefa_barramento]
command=/var/www/html/workers-control.sh tarefa_barramento
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:resposta_oficio_barramento]
command=/var/www/html/workers-control.sh resposta_oficio_barramento
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:chat]
command=/var/www/html/workers-control.sh chat
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:verificacao_virus]
command=/var/www/html/workers-control.sh verificacao_virus
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:transferir_acervo]
command=/var/www/html/workers-control.sh transferir_acervo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:assina_documento]
command=/var/www/html/workers-control.sh assina_documento
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:gestao_devedor_indice_devedor_recuperar]
command=/var/www/html/workers-control.sh gestao_devedor_indice_devedor_recuperar
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:litispendencia_manager]
command=/var/www/html/workers-control.sh litispendencia_manager
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:scheduler_cronjob_manager]
command=/var/www/html/workers-control.sh scheduler_cronjob_manager
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:judicial_processo_judicial_sincronizacao_verificacao]
command=/var/www/html/workers-control.sh judicial_processo_judicial_sincronizacao_verificacao
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:judicial_comunicacao_judicial_sincronizacao_verificacao]
command=/var/www/html/workers-control.sh judicial_comunicacao_judicial_sincronizacao_verificacao
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

#MODULO JUDICIAL:
[program:populate_comunicacao_judicial]
command=/var/www/html/workers-control.sh populate_comunicacao_judicial
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_comunicacao_judicial]
command=/var/www/html/workers-control.sh indexacao_comunicacao_judicial
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:consultar_processo]
command=/var/www/html/workers-control.sh consultar_processo
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:consultar_documento]
command=/var/www/html/workers-control.sh consultar_documento
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:consultar_avisos_pendentes]
command=/var/www/html/workers-control.sh consultar_avisos_pendentes
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:consultar_teor_comunicacao]
command=/var/www/html/workers-control.sh consultar_teor_comunicacao
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:entregar_manifestacao_processual_async]
command=/var/www/html/workers-control.sh entregar_manifestacao_processual_async
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

#MODULO RECEITA FEDERAL
[program:populate_integracao_receita_federal]
command=/var/www/html/workers-control.sh populate_integracao_receita_federal
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:ajusta_filiais_receita_federal]
command=/var/www/html/workers-control.sh ajusta_filiais_receita_federal
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_NUM_PROCS)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO